# deploy/sysvitals-backend.service
[Unit]
Description=SysVitals FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sysvitals
Group=sysvitals
WorkingDirectory=/opt/sysvitals
Environment=DATABASE_PATH=/var/lib/sysvitals/sysvitals.db
ExecStart=/opt/sysvitals/.venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
